# v2.0/deploy/cloudflared-tunnel.service
[Unit]
Description=Cloudflare Tunnel for Felswerke Terminal
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
ExecStart=/usr/bin/cloudflared tunnel --config /home/pi/.cloudflared/config.yml run
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
